# bin/noderedrevpinodes-server.service
[Unit]
Description=Revolution Pi NodeRed Service
After=multi-user.target
Documentation="https://github.com/erminas/noderedrevpinodes-server"

[Service]
Type=simple
StandardOutput=append:%L/revpi-server.log
ExecStart=/usr/bin/python3 -u revpi-server.py
WorkingDirectory=/usr/share/noderedrevpinodes-server/
Restart=on-abort

[Install]
WantedBy=multi-user.target


