# bin/noderedrevpinodes-server.service
[Unit]
Description=Revolution Pi NodeRed Service
After=multi-user.target
Documentation="https://github.com/erminas/noderedrevpinodes-server"

MemoryAccounting=true
MemoryHigh=200M
MemoryMax=400M

[Service]
Type=simple
ExecStart=/bin/sh -c 'exec /usr/bin/python3 -u revpi-server.py 2>> /var/log/revpi-server.log'
WorkingDirectory=/usr/share/noderedrevpinodes-server/
Restart=on-abort

[Install]
WantedBy=multi-user.target


